# scripts/mac/budget-sms/com.daniel.budget-sms.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.daniel.budget-sms</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/Library/Application Support/budget-sms/poll.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Application Support/budget-sms/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Application Support/budget-sms/stderr.log</string>
</dict>
</plist>
